# pi/loading-message.service
[Unit]
Description=The Loading Message — terminal client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi

# Run on tty1, fullscreen, no login prompt, no desktop.
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

Environment=TERM=linux
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 /opt/loading-message/loading_message.py

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
